# scripts/linux/aiceberg-agent.service
[Unit]
Description=AIceberg Agent
After=network-online.target
Wants=network-online.target

[Service]
User=aiceberg_agent
Group=aiceberg_agent
WorkingDirectory=/var/lib/aiceberg
EnvironmentFile=/etc/aiceberg/agent.env
ExecStart=/usr/local/bin/aiceberg_agent
Restart=always
RestartSec=3
StateDirectory=aiceberg_agent
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
